# rootfs/etc/s6-overlay/s6-rc.d/svc-spotiflac/run
#!/command/with-contenv bash

export HOME=/data
export XDG_CONFIG_HOME=/data/config
umask ${UMASK}

exec s6-setuidgid $(id -nu $PUID) /app/SpotiFLAC
